# multifolks-backend.service
[Unit]
Description=Multifolks FastAPI Backend
After=network.target mongod.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/multifolks-backend
Environment="PATH=/var/www/multifolks-backend/venv/bin"
ExecStart=/var/www/multifolks-backend/venv/bin/python /var/www/multifolks-backend/app.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=multifolks-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

